# robot/systemd/cras-robot-worker.service
[Unit]
Description=Constitutional Runtime private PiCar-X worker
After=local-fs.target

[Service]
Type=simple
User=edos
Group=edos
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/var/lib/cras-robot
WorkingDirectory=/var/lib/cras-robot
ExecStart=/usr/bin/python3 /opt/cras-robot/worker/cras_robot_worker.py
Restart=on-failure
RestartSec=2
TimeoutStopSec=3
KillSignal=SIGTERM
SendSIGKILL=yes
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cras-robot

[Install]
WantedBy=multi-user.target
